# Rinkhals/files/3-rinkhals/etc/profile.d/z_environment.sh
export TZ=UTC
export HOME=/useremain/home/rinkhals
export LD_LIBRARY_PATH=/lib:/usr/lib:$LD_LIBRARY_PATH
export PATH=/bin:/usr/bin:/usr/libexec/gcc/arm-buildroot-linux-uclibcgnueabihf/11.4.0:$PATH
export CC=/usr/bin/gcc
export CXX=/usr/bin/g++
export PYTHONPATH=
